# lib/systemd/system/var.mount
[Unit]
Requires=var-log.mount var-tmp.mount
After=var-log.mount var-tmp.mount
